# package/systemd/ercole.service
[Unit]
Description=Application

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
